# dotfiles/fcitx5.service
# ~/.config/systemd/user/fcitx5.service

[Unit]
Description=Fcitx5 Input Method
# グラフィカルセッション後に起動
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/fcitx5 -d --replace

# 入力メソッド環境変数設定
Environment=GTK_IM_MODULE=fcitx
Environment=QT_IM_MODULE=fcitx
Environment=XMODIFIERS=@im=fcitx

# 再起動設定: 失敗時に3秒待って再起動
Restart=on-failure
RestartSec=3

# リソース制限（オプション）
# MemoryMax=200M
# CPUQuota=50%

[Install]
WantedBy=default.target

# 参考:
# ステータス確認: systemctl --user status fcitx5.service
# ログ確認: journalctl --user -u fcitx5.service -f
